# ss-privoxy-client/supervisor.ini
[program:shadowsocks]
user=root
command=/usr/bin/sslocal -c /etc/shadowsocks.json
process_name=%(program_name)s
autostart=true
redirect_stderr=true

[program:privoxy]
user=root
command=/usr/sbin/privoxy --no-daemon --pidfile /run/privoxy.pid --user privoxy /etc/privoxy/config
process_name=%(program_name)s
autostart=true
redirect_stderr=true
